# opt/omd/versions/1.2.6p16.cre/skel/etc/cron.d/cmk_multisite
# Run Multisite regular jobs, e.g. scheduled reports
* * * * * curl http://localhost/###SITE###/check_mk/run_cron.py >/dev/null 2>&1
