# opt/omd/sites/Mysite/etc/cron.d/cmk_multisite
# Run Multisite regular jobs, e.g. scheduled reports
* * * * * curl http://localhost/Mysite/check_mk/run_cron.py >/dev/null 2>&1
