# opt/omd/sites/Mysite/etc/cron.d/cmk_bulk_notify
# Needed for bulk notifcations.
# Only execute cmk --notify when the microcore is currently not enabled.
* * * * * [ ! -e /omd/sites/Mysite/etc/check_mk/conf.d/microcore.mk -a -d /omd/sites/Mysite/var/check_mk/notify/bulk ] && cmk --notify send-bulks
